# contrib/cvmm@.service
[Unit]
Description=cloud-virtual-machine-manager node %i
# `modprobe@tun.service` remains for manifests that opt into `net.backend: tap`.
Wants=modprobe@tun.service modprobe@kvm.service
PartOf=machines.target
Before=machines.target
After=network.target syslog.target modprobe@tun.service modprobe@kvm.service
RequiresMountsFor=/srv/vmm/images /srv/vmm/nodes

[Service]
Type=simple
Slice=machine.slice
RestartSec=7m
TimeoutStartSec=2m
TimeoutStopSec=3m
Environment=GOMEMLIMIT=32MiB
Environment=GOGC=100
WorkingDirectory=/srv/vmm/nodes
User=hvm
Group=hvm
UMask=0077
# Default passt deployment: run the manager as a dedicated non-root service user that owns
# `<node-root>/<node>/run/` and keep `--runas` unset.
# `--runas` only changes the cloud-hypervisor child; `virtiofsd` and `passt` helpers run as this
# service user. Root manager + `--runas` is not a supported passt deployment model.
# If a manifest opts into `net.backend: tap`, add the required tap/capability handling with a
# deployment-specific override instead of reusing the default passt unit unchanged.
ExecStart=/usr/bin/cvmm \
    start %i
KillMode=mixed
Restart=on-failure

[Install]
WantedBy=machines.target